# contrib/init/soferoxd.conf
description "Soferox Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env SOFEROXD_BIN="/usr/bin/soferoxd"
env SOFEROXD_USER="soferox"
env SOFEROXD_GROUP="soferox"
env SOFEROXD_PIDDIR="/var/run/soferoxd"
# upstart can't handle variables constructed with other variables
env SOFEROXD_PIDFILE="/var/run/soferoxd/soferoxd.pid"
env SOFEROXD_CONFIGFILE="/etc/soferox/soferox.conf"
env SOFEROXD_DATADIR="/var/lib/soferoxd"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # soferoxd will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$SOFEROXD_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run soferoxd."
        echo "The setting must appear in $SOFEROXD_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following"
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is also recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"Soferox Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$SOFEROXD_PIDDIR"
    chmod 0755 "$SOFEROXD_PIDDIR"
    chown $SOFEROXD_USER:$SOFEROXD_GROUP "$SOFEROXD_PIDDIR"
    chown $SOFEROXD_USER:$SOFEROXD_GROUP "$SOFEROXD_CONFIGFILE"
    chmod 0660 "$SOFEROXD_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$SOFEROXD_PIDFILE" \
    --chuid $SOFEROXD_USER:$SOFEROXD_GROUP \
    --exec "$SOFEROXD_BIN" \
    -- \
    -pid="$SOFEROXD_PIDFILE" \
    -conf="$SOFEROXD_CONFIGFILE" \
    -datadir="$SOFEROXD_DATADIR" \
    -disablewallet \
    -daemon

